# systemd/mediaserver-autosuspend.service
[Unit]
Description=MediaServer AutoSuspend - Automatic server power management
Documentation=https://github.com/pirelike/mediaserver-autosuspend
After=network.target network-online.target
Wants=network-online.target
# Add dependencies for monitored services
After=jellyfin.service plex.service sonarr.service radarr.service nextcloud.service

[Service]
Type=simple
User=root
Group=root
# Use the virtual environment Python
ExecStart=/opt/mediaserver-autosuspend/venv/bin/python -m mediaserver_autosuspend.main
WorkingDirectory=/opt/mediaserver-autosuspend
Environment=PYTHONUNBUFFERED=1

# Configuration paths
Environment=CONFIG_PATH=/etc/mediaserver-autosuspend/config.json
Environment=LOG_PATH=/var/log/mediaserver-autosuspend

# Runtime management
Restart=on-failure
RestartSec=30
TimeoutStopSec=300

# Security settings
CapabilityBoundingSet=CAP_SYS_BOOT CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
ProtectSystem=full
ReadWritePaths=/var/log/mediaserver-autosuspend /etc/mediaserver-autosuspend
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
