# systemd/mediaserver-autosuspend.service
[Unit]
Description=MediaServer AutoSuspend Service
Documentation=https://github.com/pirelike/mediaserver-autosuspend
After=network.target jellyfin.service sonarr.service radarr.service nextcloud.service
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 -m mediaserver_autosuspend
WorkingDirectory=/opt/mediaserver-autosuspend
Environment=PYTHONPATH=/opt/mediaserver-autosuspend
Environment=PYTHONUNBUFFERED=1

# Runtime directory and permissions
RuntimeDirectory=mediaserver-autosuspend
RuntimeDirectoryMode=0755

# Configuration
EnvironmentFile=-/etc/mediaserver-autosuspend/env
ConfigurationDirectory=mediaserver-autosuspend
ConfigurationDirectoryMode=0755

# Logging
StandardOutput=append:/var/log/mediaserver-autosuspend/service.log
StandardError=append:/var/log/mediaserver-autosuspend/service.log

# Restart configuration
Restart=always
RestartSec=60

# Security settings
NoNewPrivileges=no
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true

# Resource limits
CPUQuota=25%
MemoryLimit=128M
TasksMax=100

# Capabilities
CapabilityBoundingSet=CAP_WAKE_ALARM CAP_SYS_BOOT
AmbientCapabilities=CAP_WAKE_ALARM CAP_SYS_BOOT

[Install]
WantedBy=multi-user.target
